# hw/zenoh.service
[Unit]

[Service]
ExecStart=sh -c '. /opt/ros/jazzy/setup.sh; ros2 run rmw_zenoh_cpp rmw_zenohd'

[Install]
WantedBy=default.target
